# supervisor/hongjin-es.ini
[program:hongjin-es]
command=python es.py
stdout_logfile=/server/log/hongjin-es.out
directory=/server
autostart=true
autorestart=true
startsecs=3
priority=10