# supervisor/hongjin-dbserver.ini
[program:hongjin-dbserver]
command=python dbserver.py
stdout_logfile=/server/log/hongjin-dbserver.out
directory=/server/server
autostart=true
autorestart=true
startsecs=3
priority=10