# supervisor/hongjin-cs.ini
[program:hongjin-cs]
command=python /server/server/cs.pyc
stdout_logfile=/server/server/log/hongjin-cs.out
directory=/server/server
autostart=true
autorestart=true
startsecs=3
priority=10