# supervisor/hongjin-aoiserver.ini
[program:hongjin-aoiserver]
command=python aoiserver.py
stdout_logfile=/server/log/hongjin-aoiserver.out
directory=/server
autostart=true
autorestart=true
startsecs=3
priority=10