# supervisor/hongjin-es.ini
[program:hongjin-es]
command=python /server/server/es.pyc
stdout_logfile=/server/server/log/hongjin-es.out
directory=/server/server
autostart=true
autorestart=true
startsecs=3
priority=10